# bind/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/named/
-rw-r--r--	root/root	etc/named/bind.keys
-rw-r-----	root/named	etc/named/named.conf
-rw-r-----	root/named	etc/named/rndc.key (EMPTY)
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/named
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/dig
-rwxr-xr-x	root/root	usr/bin/host
-rwxr-xr-x	root/root	usr/bin/named-checkconf
-rwxr-xr-x	root/root	usr/bin/named-checkzone
-rwxr-xr-x	root/root	usr/bin/named-compilezone
-rwxr-xr-x	root/root	usr/bin/named-journalprint
-rwxr-xr-x	root/root	usr/bin/named-rrchecker
-rwxr-xr-x	root/root	usr/bin/nsupdate
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libbind9-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libbind9.so -> libbind9-9.18.11.so
-rwxr-xr-x	root/root	usr/lib/libdns-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libdns.so -> libdns-9.18.11.so
-rwxr-xr-x	root/root	usr/lib/libirs-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libirs.so -> libirs-9.18.11.so
-rwxr-xr-x	root/root	usr/lib/libisc-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libisc.so -> libisc-9.18.11.so
-rwxr-xr-x	root/root	usr/lib/libisccc-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libisccc.so -> libisccc-9.18.11.so
-rwxr-xr-x	root/root	usr/lib/libisccfg-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libisccfg.so -> libisccfg-9.18.11.so
-rwxr-xr-x	root/root	usr/lib/libns-9.18.11.so
lrwxrwxrwx	root/root	usr/lib/libns.so -> libns-9.18.11.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/named
-rwxr-xr-x	root/root	usr/sbin/rndc
-rwxr-xr-x	root/root	usr/sbin/rndc-confgen
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/dig.1.gz
-rw-r--r--	root/root	usr/share/man/man1/host.1.gz
-rw-r--r--	root/root	usr/share/man/man1/named-checkconf.1.gz
-rw-r--r--	root/root	usr/share/man/man1/named-checkzone.1.gz
-rw-r--r--	root/root	usr/share/man/man1/named-compilezone.1.gz
-rw-r--r--	root/root	usr/share/man/man1/named-journalprint.1.gz
-rw-r--r--	root/root	usr/share/man/man1/named-rrchecker.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nsupdate.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/named.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/rndc.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/named.8.gz
-rw-r--r--	root/root	usr/share/man/man8/rndc-confgen.8.gz
-rw-r--r--	root/root	usr/share/man/man8/rndc.8.gz
drwxr-xr-x	root/root	var/
drwxrwx---	root/named	var/named/
-rw-r-----	root/named	var/named/127.in-addr.arpa
-rw-r-----	root/named	var/named/localhost
-rw-r-----	root/named	var/named/named.root
